# moinho/comun/etc/uci-defaults/91-set_libremap
#!/bin/sh
uci set libremap.settings.community='moinho'
uci set libremap.settings.community_lat='-14.065192'
uci set libremap.settings.community_lon='-47.464918'
uci commit libremap
