# moinho/comun/etc/uci-defaults/97-server-local
#!/bin/sh
uci del_list dhcp.@dnsmasq[0].address='/server.moinho.local/10.5.0.1/app'
uci add_list dhcp.@dnsmasq[0].address='/server.moinho.local/10.5.0.1/app'

