# moinho/comun/etc/uci-defaults/97_add-localserver-to-hosts
#!/bin/sh

### Adds local server do /etc/hosts

echo "10.5.0.52 www" >> /etc/hosts
